# conf/crontab/cron
# Use flock to avoid duplicate task running
0-59 * * * * /usr/bin/flock -n /tmp/crontabl.lockfile bash /period_task.sh >> /var/log/gitlab-sync/gitlab-sync.log 2>&1
